# packaging/astral-watch.service
[Unit]
Description=astral-watch — ASUS ROG Astral per-pin 12V-2x6 power monitor
Documentation=https://github.com/mbeaman/astral-watch
After=multi-user.target

[Service]
Type=simple
# '+' = run privileged (a plain ExecStartPre runs as the unprivileged User= and modprobe
# would always fail, masked by '-'); '-' = tolerate failure (i2c-dev may be built in, or
# already loaded via the shipped modules-load.d snippet)
ExecStartPre=+-/usr/sbin/modprobe i2c-dev
ExecStart=/usr/local/bin/astral-watch log /var/log/astral-watch/gpu-pins.csv --interval 0.5
Restart=always
RestartSec=5

# Least privilege: run as the unprivileged 'astral-watch' user, in group 'i2c'
# (the shipped udev rule grants that group access to the NVIDIA i2c nodes).
User=astral-watch
Group=i2c
LogsDirectory=astral-watch

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectControlGroups=true
ProtectKernelTunables=true
# AF_INET/AF_INET6 so configured ntfy/webhook notifications can leave the box
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
DevicePolicy=closed
DeviceAllow=char-i2c rw

[Install]
WantedBy=multi-user.target
